# crontab
*/2 * * * * python $NOTIFYSMS_DIR/main.py
